# crontab
0 * * * * open ~/projects/osxUtils/screen_sample.app
